# systemd/piermon.timer
[Unit]
Description=Urbit pier size monitoring

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target

